# dot_config/systemd/user/sunshine.service.d/override.conf
[Install]
# Clear the old dependency
WantedBy=
# Set the correct dependency for custom WMs
WantedBy=graphical-session.target
